# deploy/cron.d/cryptobot-scheduler
*/10 * * * * 	ubuntu	php /var/www/cryptobot/artisan schedule:run > /var/www/cryptobot/storage/logs/cron.log 2>&1
